# com.ssvcam.mjpeg-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ssvcam.mjpeg-service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/srijansrivastava/Documents/Personal/ssvcam/backend/launch_mjpeg.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/srijansrivastava/Documents/Personal/ssvcam/backend</string>
    
    <key>StandardOutPath</key>
    <string>/Users/srijansrivastava/Documents/Personal/ssvcam/backend/mjpeg_service.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/srijansrivastava/Documents/Personal/ssvcam/backend/mjpeg_service.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>MJPEG_PORT</key>
        <string>8080</string>
        <key>VIDEO_DEVICE_INDEX</key>
        <string>1</string>
    </dict>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
