# microservices/mjpeg/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:mjpeg_stream]
command=/Users/srijansrivastava/Documents/Personal/ssvcam/backend/microservices/mjpeg/venv/bin/python app.py
directory=/Users/srijansrivastava/Documents/Personal/ssvcam/backend/microservices/mjpeg
environment=PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"
autorestart=true
startretries=999
stderr_logfile=/tmp/mjpeg_stderr.log
stdout_logfile=/tmp/mjpeg_stdout.log
user=srijansrivastava
environment=PATH="/usr/local/bin:/usr/bin:/bin"